# target/linux/ipq807x/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	edgecore,eap102)
		fw_setenv upgrade_available 0
		# Unset changed flag after sysupgrade complete
		fw_setenv changed
	;;
	zyxel,nbg7815)
		# change the thermal governor policy for fan on/off
		trip_type=$(grep -ilR active /sys/class/thermal/thermal*/trip* | cut -d '/' -f 5 | tr '\n' ' ')
		for i in $trip_type; do echo bang_bang > /sys/class/thermal/$i/policy; done
	;;
	esac
}
