# target/linux/qualcommax/ipq807x/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	edgecore,eap102|\
	yuncore,ax880)
		fw_setenv upgrade_available 0
		# Unset changed flag after sysupgrade complete
		fw_setenv changed
	;;
	linksys,mx4200v1|\
	linksys,mx4200v2|\
	linksys,mx5300)
		mtd resetbc s_env || true
	;;
	esac

	#stability increase reserved mem to avoid oom when services are used that consume large amounts of mem e.g. samba with attached usb disk etc.
	sysctl -w vm.min_free_kbytes=65536
}
